# crontab/user_crontab
* 12 * * * ~/clippings_backup_autoremove.sh
